# systemd/hermes-agent.service
[Unit]
Description=Hermes Agent Dashboard
After=network.target

[Service]
Type=simple
User=hermes
ExecStart=/home/hermes/.hermes/hermes-agent/venv/bin/hermes dashboard --port 9119 --host 127.0.0.1
Restart=always
RestartSec=5
ProtectSystem=strict
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
