# configs/oauth2-proxy.service
[Unit]
Description=OAuth2 Proxy
After=network.target

[Service]
Type=simple
User=oauth2-proxy
EnvironmentFile=/opt/secrets/hermes.env
ExecStart=/usr/local/bin/oauth2-proxy \
  --provider=oidc \
  --oidc-issuer-url=https://oauth.yandex.ru \
  --login-url=https://oauth.yandex.ru/authorize \
  --redeem-url=https://oauth.yandex.ru/token \
  --profile-url=https://login.yandex.ru/info \
  --scope="login:email login:info" \
  --client-id=${YANDEX_CLIENT_ID} \
  --client-secret=${YANDEX_CLIENT_SECRET} \
  --cookie-secret=${OAUTH2_PROXY_COOKIE_SECRET} \
  --cookie-secure=true \
  --cookie-domain=ai.mais.agency \
  --upstream=http://127.0.0.1:8080 \
  --http-address=127.0.0.1:4180 \
  --set-xauthrequest=true \
  --pass-access-token=false \
  --email-domain=* \
  --authenticated-emails-file=/etc/oauth2-proxy/authenticated-emails.txt

Restart=always
RestartSec=5
ProtectSystem=strict
PrivateTmp=yes
NoNewPrivileges=yes
PrivateDevices=yes

[Install]
WantedBy=multi-user.target
